# recipes-myapps/pin-monitor/files/pin-config.service
[Unit]
Description=TRT Pin Configuration - Configuración inicial de pines
Documentation=https://github.com/AntonioTRT/TRT_yocto
Before=pin-monitor.service
After=local-fs.target

[Service]
Type=oneshot
User=root
Group=root
ExecStart=/usr/bin/python3 /usr/bin/pin-config.py
ExecStartPost=/bin/chown -R trt:trt /opt/trt/config
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal
SyslogIdentifier=trt-pin-config

# Timeout
TimeoutStartSec=30

[Install]
WantedBy=multi-user.target
RequiredBy=pin-monitor.service